# installers/launchd/com.magpi.indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.magpi.indexer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/magpi</string>
        <string>index</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/magpi-indexer.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/magpi-indexer.err</string>
</dict>
</plist>
